# assets/supervisord.conf
[supervisord]
nodaemon=true

[program:oracle-xe]
command=/usr/sbin/startup.sh

[program:oracle-tail]
command=/usr/bin/tail -F /u01/app/oracle/diag/tnslsnr/${HOSTNAME}/listener/trace/listener.log
